# cdr_transform/cdr_transform.cron
# CDR Transform - Process CDR files every hour
# Place this file in /etc/cron.d/cdr_transform

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# CDR processing - runs at :00 every hour
0 * * * * {USER} python3 {PROJECT_DIR}/cdr_transform/cdr_transform.py cdr

# LU processing - runs at :30 every hour
30 * * * * {USER} python3 {PROJECT_DIR}/cdr_transform/cdr_transform.py lu
